# lego-wrapper.cron.weekly
#!/bin/sh

. /etc/lego-wrapper.d/lego-wrapper.conf

$LEGO_WRAPPER -R
